# fanctl.service
[Unit]
Description=Fan control based on cpu temperature

[Service]
Type=simple
ExecStart=/usr/local/bin/fanctl

[Install]
WantedBy=multi-user.target
